# config/com.sovi.dashboard.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.sovi.dashboard</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/noh/Work/ai/sovi/.venv/bin/python</string>
        <string>-m</string>
        <string>sovi</string>
        <string>server</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/noh/Work/ai/sovi</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PYTHONPATH</key>
        <string>/Users/noh/Work/ai/sovi/src</string>
    </dict>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/noh/Work/ai/sovi/output/dashboard.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/noh/Work/ai/sovi/output/dashboard.err</string>
    <key>RunAtLoad</key>
    <true/>
</dict>
</plist>
